# runit/pi-sc-lcd-remind/run
#!/bin/sh
export TCLLIBPATH="$TCLLIBPATH /home/pi/sc/bin"
exec chpst -u pi:pi expect /home/pi/sc/bin/lcdproc-remind.expect >/dev/null
