# runit/pi-sc-doormon/run
#!/bin/sh
exec /home/pi/sc/bin/door-monitor.sh
