# runit/pi-sc-remind/run
#!/bin/sh

exec chpst -u pi:pi:dialout remind -z /home/pi/sc/remind
