# runit/pi-sc-lcd/run
#!/bin/sh
exec chpst -u pi:pi expect /home/pi/sc/bin/lcdproc.expect >/dev/null
