# runit/pi-sc-lcd/run
#!/bin/sh
export TCLLIBPATH="$TCLLIBPATH /home/pi/sc/bin"
exec chpst -d 20000000 -u pi:pi expect /home/pi/sc/bin/lcdproc-monitor.expect >/dev/null
