# runit/pi-sc-monitor/run
#!/bin/sh
exec chpst -u pi:pi:dialout:i2c python -u /home/pi/sc/bin/pidmonitor.py
