# runit/pi-sc-monitor/run
#!/bin/sh
exec chpst -u pi:pi:dialout python -u /home/pi/sc/bin/monitor.py
