# runit/pi-sc-monitor280/run
#!/bin/sh
exec chpst -u pi:pi:i2c python -u /home/pi/sc/bin/monitor280.py
